# scripts/cardirehab-backend-preprod.service
[Unit]
Description=CardiRehab Backend Spring Boot Application (Preprod)
After=network.target mariadb.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/html/backend-preprod
ExecStart=/usr/lib/jvm/java-17-openjdk-amd64/bin/java $JAVA_OPTS -jar /var/www/html/backend-preprod/target/crehab.war
ExecStop=/bin/kill -15 $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cardirehab-backend-preprod

EnvironmentFile=/var/www/html/backend-preprod/.env.preprod
Environment="SPRING_PROFILES_ACTIVE=preprod"
Environment="JAVA_OPTS=-Xms512m -Xmx1024m"

[Install]
WantedBy=multi-user.target
