# scripts/cardirehab-backend-prod.service
[Unit]
Description=CardiRehab Backend Spring Boot Application (Production)
After=network.target mariadb.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/html/backend
ExecStart=/usr/lib/jvm/java-17-openjdk-amd64/bin/java $JAVA_OPTS -jar /var/www/html/backend/target/crehab.war
ExecStop=/bin/kill -15 $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cardirehab-backend-prod

EnvironmentFile=/var/www/html/backend/.env.prod
Environment="SPRING_PROFILES_ACTIVE=prod"
Environment="JAVA_OPTS=-Xms512m -Xmx1024m"

[Install]
WantedBy=multi-user.target
